# pLUM.timer
[Unit]
Description=Daily pLUM

[Timer]
OnCalendar=*-*-* %HOUR%
Persistent=true

[Install]
WantedBy=timers.target
