# pLUM-package/DEBIAN/prerm
#!/usr/bin/env bash

echo Stopping and disabling pLUM.timer ...

/usr/bin/systemctl disable --now pLUM.timer
